# systemd/relay-daemon.service
# relay-daemon.service — systemd unit for the relay-tools HTTP API daemon
#
# Installation
# ------------
# 1. Copy this file to the systemd system directory:
#       sudo cp relay-daemon.service /etc/systemd/system/
#
# 2. Copy the sample configs (edit to taste):
#       sudo mkdir -p /etc/relay
#       sudo cp channels.yaml /etc/relay/channels.yaml
#       sudo mkdir -p /etc/relay/boards.d
#       sudo cp rom2820-board.yaml /etc/relay/boards.d/rom2820.yaml
#
# 3. Reload the systemd manager:
#       sudo systemctl daemon-reload
#
# 4. Enable the service to start automatically on boot:
#       sudo systemctl enable relay-daemon.service
#
# 5. Start the service immediately:
#       sudo systemctl start relay-daemon.service
#
# 6. Check its status:
#       sudo systemctl status relay-daemon.service
#
# Configuration
# -------------
# Channel initial states are read from /etc/relay/channels.yaml on every
# start.  All channels default to OFF if the file is absent or a channel
# is not listed.  See channels.yaml for the file format.
#
# Override any setting without editing this file by using a drop-in:
#       sudo systemctl edit relay-daemon.service
#
# Common overrides:
#   - Change listen port:        ExecStart=/usr/bin/relay serve --port 9000
#   - Force a specific driver:   Environment=RELAY_DRIVER=rpigpio
#   - Bind to localhost only:    ExecStart=… --host 127.0.0.1
#   - Use a different config:    Environment=RELAY_CONFIG=/etc/relay/custom.yaml
#   - Board CLI profile path:    export RELAY_BOARD_CONFIG=/etc/relay/boards.d/rom2820.yaml
#
# Requirements
# ------------
# relay-tools must be installed (with GPIO extras) before enabling this unit:
#   pip install relay-tools[gpio]
#
# On Raspberry Pi the process user must be in the 'gpio' group, or you can
# run as root (SupplementaryGroups= / User= below).

[Unit]
Description=relay-tools HTTP API daemon
Documentation=https://github.com/Advantech-EECC/relay-tools
After=network.target
Wants=network.target

[Service]
Type=simple

# Run as root so GPIO access is available without extra group configuration.
# To run as a dedicated user instead, create that user, add it to the 'gpio'
# group, and change the two lines below:
#   User=relay
#   SupplementaryGroups=gpio
User=root

# Ensure the config directory exists before starting.
ExecStartPre=-/bin/mkdir -p /etc/relay
ExecStartPre=-/bin/mkdir -p /etc/relay/boards.d

# Path to the relay CLI entry point installed by pip.
# Adjust if you installed into a virtualenv, e.g.:
#   /opt/relay-tools/bin/relay serve --host 0.0.0.0 --port 8000
ExecStart=/usr/bin/relay serve --host 0.0.0.0 --port 8000

# GPIO driver: "auto" (default) tries RPi.GPIO first, then gpiozero.
# Override with "rpigpio" or "gpiozero" if auto-detection is unreliable.
Environment=RELAY_DRIVER=auto

# Path to the per-channel initial-state config file.  All channels start
# OFF if the file is absent; channels not listed also default to OFF.
Environment=RELAY_CONFIG=/etc/relay/channels.yaml

# Optional default profile path for relay-board operators:
# Environment=RELAY_BOARD_CONFIG=/etc/relay/boards.d/rom2820.yaml

# Restart the daemon automatically if it crashes.
Restart=on-failure
RestartSec=5s

# Give the process time to shut down cleanly before SIGKILL.
TimeoutStopSec=10s

# Write stdout/stderr to the systemd journal (view with journalctl -u relay-daemon).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-daemon

[Install]
WantedBy=multi-user.target
